# wallpaper.service
[Unit]
Description=Wallpaper

[Service]
ExecStart=/filepath/wallpaper

[Install]
WantedBy=default.target